# image-overlay/etc/rc.local
# echo sAaAj > /dev/led/beaglebone:green:heartbeat
